# systemd/snd-repair-fw-validation-suspend.service
[Unit]
Description=snd_repair FW validation after suspend/resume (TAS2783 / PX13)
Documentation=file:@REPO_ROOT@/validation/README.md
# Mismo grafo que px13-audio-resume.service (post-resume, no en hook system-sleep)
After=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
After=px13-audio-resume.service

[Service]
Type=oneshot
WorkingDirectory=@REPO_ROOT@
Environment=SND_REPAIR_REPO=@REPO_ROOT@
Environment=SND_REPAIR_SUSPEND_DELAY=15
ExecStart=@REPO_ROOT@/scripts/fw-validation-suspend-hook.sh

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
